# contrib/systemd/xdg-desktop-portal-hyprland.service.in
[Unit]
Description=Portal service (Hyprland implementation)
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.hyprland
ExecStart=@libexecdir@/xdg-desktop-portal-hyprland
Restart=on-failure
