# contrib/systemd/xdg-desktop-portal-wlr.service.in
[Unit]
Description=Portal service (wlroots implementation)

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.wlr
ExecStart=@libexecdir@/xdg-desktop-portal-wlr
Restart=on-failure
